# roles/podluck_pod_systemd/files/systemd-skel/podluck.target
[Unit]
Description=Podman pod %p containers managed by podluck
PartOf=%p.service
After=%p.service
